# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/fd/1
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/app/docker/tee-log.sh postgresql/postgresql /usr/lib/postgresql/18/bin/postgres -D /var/lib/postgresql/data -c logging_collector=off -c log_destination=stderr
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
startsecs=5
stopasgroup=true
killasgroup=true

[program:app]
directory=/app
command=python /app/run.py
user=dbclaw
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
startsecs=10
